# forge/dependabot/dependabot-bumper.service
[Unit]
Description=Dependabot bumper - scan, gate, and auto-merge one clean low-risk dependency bump
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# TEMPLATE: this in-repo unit uses neutral placeholder paths. The installed copy
# (e.g. ~/.config/systemd/user/dependabot-bumper.service) is machine-local — edit
# THAT copy with your real checkout path and PATH; do not point systemd at this file.
# %h expands to the unit user's home; adjust "%h/forge" to wherever the repo lives.
ExecStart=%h/forge/forge/dependabot/run.sh
WorkingDirectory=%h/forge
# Machine-local values (router endpoint, API key, dirs) live in a gitignored .env at
# the repo root; it is plain KEY=VALUE, so it doubles as a systemd EnvironmentFile.
# EnvironmentFile=%h/forge/.env

Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

# Suite + audit + sign-off can take a while on a busy day.
TimeoutStartSec=2400

# Exit codes: 0 = merged/branched/no-candidates; 1 = advisory or error. An advisory run
# therefore shows this unit as "failed" — that is deliberate visibility: the details are in
# the Forge task it filed and in forge/dependabot/logs/auto.jsonl.

# NO sandboxing, deliberately (first supervised run finding, 2026-07-06): in USER units,
# systemd implements ProtectSystem/PrivateTmp/NoNewPrivileges via unprivileged user
# namespaces, which remap root-owned files to nobody — breaking both sudo (gaol's watchdog
# under run_tests: "/etc/sudo.conf owned by uid 65534") and ssh pushes ("bad owner on
# /etc/ssh/ssh_config.d"). This is trusted first-party automation running as the user; the
# bumper's own gates are the containment.

[Install]
WantedBy=default.target
