# forge/dependabot/dependabot-bumper.service
[Unit]
Description=Dependabot bumper - scan, gate, and auto-merge one clean low-risk dependency bump
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/user/Projects/erewhon/forge/forge/dependabot/run.sh
WorkingDirectory=/home/user/Projects/erewhon/forge

Environment=HOME=/home/user
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/home/user/.local/bin:/usr/local/bin:/usr/bin:/bin

# Suite + audit + sign-off can take a while on a busy day.
TimeoutStartSec=2400

# Exit codes: 0 = merged/branched/no-candidates; 1 = advisory or error. An advisory run
# therefore shows this unit as "failed" — that is deliberate visibility: the details are in
# the Forge task it filed and in forge/dependabot/logs/auto.jsonl.

# NO sandboxing, deliberately (first supervised run finding, 2026-07-06): in USER units,
# systemd implements ProtectSystem/PrivateTmp/NoNewPrivileges via unprivileged user
# namespaces, which remap root-owned files to nobody — breaking both sudo (gaol's watchdog
# under run_tests: "/etc/sudo.conf owned by uid 65534") and ssh pushes ("bad owner on
# /etc/ssh/ssh_config.d"). This is trusted first-party automation running as the user; the
# bumper's own gates are the containment.

[Install]
WantedBy=default.target
